# pkg/deb/portfwd_V.SS-R_i386/DEBIAN/postinst
#! /bin/sh

update-rc.d portfwd.init defaults >/dev/null
