# pkg/deb/portfwd_V.SS-R_i386/DEBIAN/postrm
#! /bin/sh

portfwd_init=/etc/init.d/portfwd.init

[ -x "$portfwd_init" ] && $portfwd_init stop

if [ purge = "$1" ]; then
	update-rc.d $portfwd_init remove >/dev/null
fi
